# install/debian/postrm
#! /bin/sh

case "$1" in
    remove)
        update-rc.d prax remove
        #sed -i 's/ prax//' /etc/nsswitch.conf
        ;;
    purge)
        ;;
esac
